# pihole/services/padd/finish
#!/usr/bin/with-contenv bash

s6-echo "Stopping PADD"
killall -9 padd.sh
