# pihole/services/padd/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

s6-echo "Starting PADD"
sleep 30
/usr/src/padd.sh 2> /dev/null > /dev/tty1
