# pihole/services/fbcp/finish
#!/usr/bin/with-contenv bash

s6-echo "Stopping fbcp"
killall -9 fbcp
