# pihole/services/fbcp/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

s6-echo "Starting fbcp"
/usr/src/fbcp
